# .config/systemd/user/gnome-keyring-daemon.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/gnome-keyring-daemon --foreground --components=pkcs11,secrets --control-directory=%t/keyring --login
